# config/etc/inittab
::sysinit:/etc/init.d/rcS

tty1::respawn:/sbin/getty -n -l /bin/ash 38400 tty1 vt100
tty2::respawn:/sbin/getty -n -l /bin/ash 38400 tty2 vt100
tty3::respawn:/sbin/getty -n -l /bin/ash 38400 tty3 vt100
tty4::respawn:/sbin/getty -n -l /bin/ash 38400 tty4 vt100

::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
::restart:/sbin/init
